# getquotes.cron
0 10 * * * /cygdrive/c/dev/repos/StockApp/getquotes.sh
